# packages/datapool/datapool-0.3.0.tar.gz/datapool-0.3.0/scripts/datapool.service
[Unit]
Description=eawag datapool service

[Service]
WorkingDirectory=/
Environment=HOME=/root
ExecStart=/usr/local/bin/pool run-simple-server
Restart=always

[Install]
WantedBy=multi-user.target
